# deploy/systemd/batch-eod.service
[Unit]
Description=AEQCS end-of-day batch

[Service]
Type=oneshot
WorkingDirectory=/opt/aeqcs
ExecStart=/opt/aeqcs/.venv/bin/python -m aeqcs.runtime.batch eod
