# deploy/systemd/batch-eod.timer
[Unit]
Description=Run AEQCS end-of-day batch

[Timer]
OnCalendar=Mon..Fri 15:10
Persistent=true

[Install]
WantedBy=timers.target
