# deploy/systemd/db-vacuum.timer
[Unit]
Description=Run AEQCS PostgreSQL vacuum maintenance nightly

[Timer]
OnCalendar=Mon..Fri 23:40:00 Asia/Shanghai
Persistent=true

[Install]
WantedBy=timers.target
