# deploy/systemd/intraday.service
[Unit]
Description=AEQCS intraday monitor

[Service]
Type=simple
User=aeqcs
EnvironmentFile=-/etc/aeqcs/aeqcs.env
WorkingDirectory=/opt/aeqcs
MemoryMax=16G
ExecStart=/opt/aeqcs/.venv/bin/python -m aeqcs.runtime.intraday
RuntimeMaxSec=21000
Restart=on-failure
