# deploy/systemd/intraday.timer
[Unit]
Description=Start AEQCS intraday monitor

[Timer]
OnCalendar=Mon..Fri 09:15
Persistent=true

[Install]
WantedBy=timers.target
