# deploy/systemd/mcp-server.service
[Unit]
Description=AEQCS deterministic core MCP HTTP/SSE server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=aeqcs
EnvironmentFile=-/etc/aeqcs/aeqcs.env
WorkingDirectory=/opt/aeqcs
MemoryMax=16G
ExecStart=/opt/aeqcs/.venv/bin/python -m aeqcs.core.mcp_server
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aeqcs/data /opt/aeqcs/logs

[Install]
WantedBy=multi-user.target
