# deploy/systemd/restore-rehearsal.service
[Unit]
Description=AEQCS restore rehearsal

[Service]
Type=oneshot
User=aeqcs
EnvironmentFile=-/etc/aeqcs/aeqcs.env
WorkingDirectory=/opt/aeqcs
MemoryMax=16G
ExecStart=/opt/aeqcs/.venv/bin/python -m aeqcs.runtime.batch restore-rehearsal
