# deploy/systemd/restore-rehearsal.timer
[Unit]
Description=Run AEQCS restore rehearsal weekly

[Timer]
OnCalendar=Sun 02:30:00 Asia/Shanghai
Persistent=true

[Install]
WantedBy=timers.target
